# cookbooks/lb_aicache/templates/default/aicache_restart.sh.erb
#!/bin/bash
# aicache daemon
# chkconfig: 345 20 80
# description: aicache daemon
# processname: aicache

DAEMON_PATH="/usr/local/aicache"
RELOAD_FILE="/usr/local/aicache/reload"
RELOAD_SUCCESS="/usr/local/aicache/reload_success"
RELOAD_FAIL="/usr/local/aicache/reload_fail"

DAEMON=aicache
DAEMON_FULL=$DAEMON_PATH/$DAEMON
DAEMONOPTS="-f /etc/aicache/aicache.cfg -c rs"

NAME=aicache
DESC="aiCache - Caching and DoS protection solution"
PIDFILE=/var/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME

case "$1" in
start)
        printf "%-50s" "Starting $NAME..."
        cd $DAEMON_PATH
        export LD_LIBRARY_PATH=/usr/local/lib
        echo 256000 > /proc/sys/fs/file-max        # Adjust system-wide number of open files
        echo 4 > /proc/sys/net/ipv4/tcp_fin_timeout
        ulimit -c 104857600
        ulimit -n 65535
        $DAEMON_FULL $DAEMONOPTS
        nohup /usr/local/aicache/alert.pl > /dev/null 2>&1 &
        PID=`pgrep -f $DAEMON_FULL`
        if [ -z $PID ]; then
            printf "%s\n" "Fail"
        else
            echo $PID > $PIDFILE
            printf "%s\n" "Ok"
        fi
;;
status)
        printf "%-50s" "Checking $NAME..."
        if [ -f $PIDFILE ]; then
            PID=`cat $PIDFILE`
            if [ -z "`ps axf | grep ${PID} | grep -v grep`" ]; then
                printf "%s\n" "Process dead but pidfile exists"
            else
                echo "Running"
            fi
        else
            printf "%s\n" "Service not running"
        fi
;;
stop)
        printf "%-50s" "Stopping $NAME"
            PID=`cat $PIDFILE`
            cd $DAEMON_PATH
        if [ -f $PIDFILE ]; then
            kill -HUP $PID
            printf "%s\n" "Ok"
            rm -f $PIDFILE
            pkill -f /usr/local/aicache/alert.pl > /dev/null 2>&1
        else
            printf "%s\n" "pidfile not found"
        fi
;;

restart)
        $0 stop
        $0 start
;;

reload)
        if [ -z `pidof aicache` ]; then
                echo "aiache not running ... starting"
                $0 start
        else
                echo "aicache running ... reloading"
                printf "%-50s" "Reloading $NAME"
                touch $RELOAD_FILE
                if [ -f $RELOAD_SUCCESS ]; then
                        printf "%s\n" "Ok"
                else
                        printf "%s\n" "there was an error"
                        cat $RELOAD_FAIL
                fi
        fi

;;

*)
        echo "Usage: $0 {status|reload|start|stop|restart}"
        exit 1
esac
